# crontab
# START CRON JOB
* * * * * python3 /app/extract.py
# END CRON JOB
